# docker/mysql_sh/supervisord.conf
[supervisord]
nodaemon=true
logfile=/home/rei/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=false

[unix_http_server]
file=/home/rei/supervisor/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/rei/supervisor/supervisor.sock



# [program:mysqld]
# command=gosu mysql mysqld
# autostart=true
# autorestart=true
# stdout_logfile=/var/log/mysqld.log

[program:sqlinit]
command=/bin/bash -c "sleep 16 && cd /home/rei/mysql_shell && ./sqlinit"
stdout_logfile=/home/rei/mysql_shell/sqlinit.log
stderr_logfile=/home/rei/mysql_shell/sqlinit.err
autostart=true
autorestart=false
startsecs=0

[program:rei]
command=/bin/bash -c "sleep 20 && cd /home/rei/sunny_peace_rei && /home/rei/sunny_peace_rei/rei -config /home/rei/sunny_peace_rei/conf/config.yaml"
stdout_logfile=/home/rei/sunny_peace_rei/supervisor/rei_app.log
stderr_logfile=/home/rei/sunny_peace_rei/supervisor/rei_app.err
autostart=true
autorestart=false
startsecs=2
startretries=2



